# roles/restic/templates/restic-backup.timer.j2
[Unit]
Description={{ restic_backup_name }} backup timer

[Timer]
OnCalendar={{ restic_schedule }}
Unit=restic-{{ restic_backup_name }}.service

[Install]
WantedBy=timers.target
