# roles/torrent/files/flexget.timer
[Unit]
Description=Run flexget every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
