# roles/torrent/files/flexget.service
[Unit]
Description=Flexget Daemon

[Service]
Type=simple
User=torrent
Group=torrent
WorkingDirectory=/home/torrent/flexget
ExecStart=/usr/local/bin/flexget daemon start
ExecStop=/usr/local/bin/flexget daemon stop
ExecReload=/usr/local/bin/flexget daemon reload

[Install]
WantedBy=multi-user.target
