# roles/torrent/files/wg-quick@.service
[Service]
NetworkNamespacePath=/run/netns/torrent
BindPaths=/etc/netns/torrent/resolv.conf:/etc/resolv.conf
