# roles/share/files/photoprism-index.timer
[Unit]
Description=Schedules photoprism index

[Timer]
OnCalendar=*-*-* 03:30:00
Unit=photoprism-index.service

[Install]
WantedBy=timers.target
