# roles/duckdns/files/duckdns.timer
[Unit]
Description=Run DuckDNS every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
